# iot_1.2.0_kernel_3.19.8/meta-intel-galileo/recipes-kernel/galileo-init/files/galileo-init.service
[Unit]
Description=Load modules for the Intel galileo
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/sh /usr/sbin/galileo-init.sh

[Install]
WantedBy=local-fs.target
